# systemd/speedreport.service
[Unit]
Description=Ookla Speedtest Report

[Service]
ExecStart=/home/pi/projects/speedreport/systemd/report.sh
